# com.valor.autoexperiment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.valor.autoexperiment</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>__PROJECT_DIR__/scripts/autoexperiment.py</string>
        <string>--target</string>
        <string>__TARGET__</string>
        <string>--iterations</string>
        <string>100</string>
        <string>--budget</string>
        <string>2.0</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/autoexperiment.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/autoexperiment.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>__PROJECT_DIR__</string>
    </dict>
</dict>
</plist>
